# agent/app/cabsagent.service
#!/bin/bash
[Unit]
Description="Connection Automation Brokerage System (CABS)"

[Service]
ExecStart=/opt/cabsagent/cabsagent.py

[Install]
WantedBy=multi-user.target
