# graph/conf/cabsgraph.service
#!/bin/bash
[Unit]
Description="CABS graphing"

[Service]
ExecStart=/opt/cabsgraph/update.sh

[Install]
WantedBy=multi-user.target
